# flask_cap/templates/etc/supervisor_celery.conf
[program:celery_mail]
command=/root/v3/bin/celery worker -A celery_worker.celery -l info -Q task_email -n worker_email@%h
directory=/root/{{ project_name }}
user=root
autostart=true
autorestart=true
stdout_logfile=/root/logs/celery_mail.log
stderr_logfile=/root/logs/celery_mail.log
